# deploy/systemd/denotary-audit.service
[Unit]
Description=deNotary Audit API
After=network-online.target denotary-finality-watcher.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/denotary/offchain.env
WorkingDirectory=/
ExecStart=/usr/bin/env bash -lc 'exec "$PROJECT_ROOT/scripts/run-offchain-service.sh" audit'
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
